# matlablm.service
[Unit]
Description=MatlabFlexLM
Documentation=https://github.com/jstaf/matlablm-systemd
After=network.target

[Service]
User=CHANGE_ME
Type=forking
ExecStart=/usr/local/MATLAB/R2016b/etc/lmstart
ExecStop=/usr/local/MATLAB/R2016b/etc/lmdown

[Install]
WantedBy=multi-user.target
